# src/alfred/_bundled/systemd/alfred-stayc-clinical.service.template
[Unit]
Description=STAY-C — sovereign ambient clinical scribe (PHI · local-only · no cloud egress)
After=network.target
# DELIBERATELY no PartOf= / no WantedBy=algernon.target — standalone lifecycle,
# never joins the salem/kal-le/hypatia fan-out (GROUND #1,#6). STAY-C is NOT a
# row in ~/.alfred/instances.yaml; it is installed by its OWN dedicated
# installer (alfred.scripts.install_stayc_unit), never by install_systemd_units.

[Service]
Type=simple
# SYSTEM unit (#67): installed to /etc/systemd/system + enabled by root, so a
# PRIVILEGED system manager applies the hardening below (capability drops,
# namespaces, ProtectSystem=strict) THEN drops to this unprivileged User/Group
# to run the scribe. A `systemctl --user` manager is UNPRIVILEGED and CANNOT
# apply that hardening — it dies status=218/CAPABILITIES ("Failed to drop
# capabilities", "UTS namespace setup prohibited") and crash-loops, never
# booting. That is why STAY-C is a system unit, not a --user unit.
User=<STAYC_USER>
Group=<STAYC_GROUP>
# WorkingDirectory is the STAY-C data root; the ExecStart python is the STAY-C
# OWN venv (stayc-root/.venv) — faster-whisper lives there, NOT in the shared
# alfred venv (#67 F1). `python -m alfred` resolves alfred from that venv's
# installed/editable package, so WorkingDirectory need NOT be the code repo.
WorkingDirectory=<STAYC_WORKDIR>
ExecStart=<STAYC_PYTHON> -m alfred --config <STAYC_CONFIG> up --_internal-foreground
# --- barrier-(c) cloud-key scrub (DETERMINATE B): manager-enforced regardless of launch, supersedes the header's env -u wrapper ---
UnsetEnvironment=ANTHROPIC_API_KEY ANTHROPIC_AUTH_TOKEN ANTHROPIC_BASE_URL GROQ_API_KEY DEEPGRAM_API_KEY OPENAI_API_KEY OPENROUTER_API_KEY ELEVENLABS_API_KEY TOGETHER_API_KEY ZO_API_KEY RESEND_API_KEY
# salt-only, 0600, sited OUTSIDE the config dir so barrier-(c)'s sibling-.env autoload never sees it; carries NO cloud key. It sources SCRIBE_ENCOUNTER_SALT, which config.scribe.encounter_salt references as ${SCRIBE_ENCOUNTER_SALT} (env-sourced, NEVER inline).
# Leading `-` (#67 F2): OPTIONAL — a MISSING env-file is non-fatal to STARTUP (without the dash systemd would REFUSE to start the unit). Fail-closed is preserved at INGEST, not startup: if the file is absent the salt stays the unresolved literal ${SCRIBE_ENCOUNTER_SALT}, and compute_encounter_id (scribe/identity.py) REJECTS an unresolved-placeholder salt — so the daemon boots idle-ready but refuses to mint a reversible encounter_id until the salt is set. Startup-tolerant + ingest-fail-closed, never a weak id.
EnvironmentFile=-<STAYC_SECRETS_ENV>
# real-need iii: relocate HF STT cache off /home + force offline model load
Environment=HF_HOME=<STAYC_HF_HOME>
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
# ReadWritePaths carves the writable holes out of ProtectSystem=strict: vault +
# data + HF cache, PLUS the scribe input_dir (#67 F4). input_dir is derived
# from config.scribe.input_dir by the installer — on-box it is a SIBLING of the
# data dir, NOT under it, so ProtectSystem=strict would leave it read-only and
# every encounter would EROFS on its ledger write. Include it explicitly.
ReadWritePaths=<STAYC_VAULT> <STAYC_DATA> <STAYC_HF_HOME> <STAYC_INPUT_DIR>
# <STAYC_SEAL_DIR> holds the seal PUBLIC key + the s.50 retention schedule — both
# OPERATOR-published (via `alfred scribe retention keygen` / `schedule publish`,
# run in the operator's shell, NOT under this sandbox) and daemon-READ-ONLY (#13d
# Q4). Siting them read-only means a live-compromised daemon can seal TO the key
# but cannot SWAP it (seal future encounters to an attacker key) or rewrite the
# schedule. The sealed-blob STORE stays under <STAYC_DATA> (ReadWritePaths) — the
# daemon writes blobs there; only the pubkey + schedule are read-only.
ReadOnlyPaths=<STAYC_CONFIG> <STAYC_SEAL_DIR>
PrivateTmp=yes
UMask=0077
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# best-effort egress belt — PROBED at boot, LOGGED enforced|unverified, NON-GATING (Approach B; GROUND #8). Order irrelevant (both accumulate); allow keeps 127.0.0.0/8 + ::1 reachable even when enforced.
IPAddressAllow=localhost
IPAddressDeny=any
Restart=on-failure
RestartSec=30s
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStopSec=30s
# 78 = missing-deps sentinel; 79 = sovereign boundary breach — load-bearing ONLY
# because of the orchestrator's sovereign-gated exit-79 propagation (GROUND #7):
# the --_internal-foreground plain monitor loop tears every child down and
# re-raises 79 so this directive keeps STAY-C DOWN. (The live TUI/dashboard
# supervision paths do NOT yet honor 79 — task #59; STAY-C runs non-live so it
# is unaffected today.)
RestartPreventExitStatus=78 79
# DELIBERATELY OMITTED (would break real-need iii — STT native inference):
# MemoryDenyWriteExecute (kills ctranslate2/torch W+X JIT) and
# SystemCallFilter=@system-service (risks silently killing native STT). Left as
# an on-box-validated future tightening, NOT in the shipped load-bearing set.

[Install]
WantedBy=multi-user.target
